# packaging/packetwyrm-proxyd.service
[Unit]
Description=PacketWyrm Web GUI / remote-RPC gateway (HTTPS -> control socket)
Documentation=https://github.com/kazubu/packetwyrm
After=packetwyrmd.service
Wants=packetwyrmd.service

[Service]
Type=simple
# TLS-terminating relay to the daemon's Unix control socket. Serves the
# embedded Web GUI and POST /api/rpc. Self-signed cert by default; pass
# --tls-cert/--tls-key for a real certificate. See docs/design/web-gui.md.
# Binds 127.0.0.1 by default so a fresh install starts cleanly (a public
# 0.0.0.0 bind fail-closes unless system.secret is set -- which the seeded
# example config isn't). To expose it, set system.secret in the env config
# and change --listen to 0.0.0.0:8443 (deliberate opt-in).
ExecStart=/usr/local/sbin/packetwyrm-proxyd \
    --listen 127.0.0.1:8443 \
    --socket /run/packetwyrm/packetwyrmd.sock
Restart=on-failure
RestartSec=2s

# Network-facing and holds no hardware: runs unprivileged as the packetwyrm
# user/group (created by packaging/packetwyrm.sysusers). It reaches the daemon's
# control socket, which the daemon creates 0660 root:packetwyrm in production --
# so this gateway connects via its packetwyrm group membership, not as root. The
# daemon remains the sole authority for the access secret (forwarded verbatim).
User=packetwyrm
Group=packetwyrm

# Hardening.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectControlGroups=true
ProtectKernelTunables=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
# No writable paths needed (self-signed cert is generated in memory).

[Install]
WantedBy=multi-user.target
