# packaging/packetwyrm-proxyd.service
[Unit]
Description=PacketWyrm Web GUI / remote-RPC gateway (HTTPS -> control socket)
Documentation=https://github.com/kazubu/packetwyrm
After=packetwyrmd.service
Wants=packetwyrmd.service

[Service]
Type=simple
# TLS-terminating relay to the daemon's Unix control socket. Serves the
# embedded Web GUI and POST /api/rpc. All options live in the config file
# /etc/packetwyrm/proxyd.yaml (listen / socket / tls_cert / tls_key / no_tls /
# insecure_no_auth / allowed_hosts) -- edit that and `systemctl restart
# packetwyrm-proxyd`, no unit override needed. It ships loopback-only + self-
# signed TLS; a public 0.0.0.0 bind fail-closes unless system.secret is set or
# insecure_no_auth is enabled. See docs/design/web-gui.md.
ExecStart=/usr/local/sbin/packetwyrm-proxyd --config /etc/packetwyrm/proxyd.yaml
Restart=on-failure
RestartSec=2s

# Network-facing and holds no hardware: runs unprivileged as the packetwyrm
# user/group (created by packaging/packetwyrm.sysusers). It reaches the daemon's
# control socket, which the daemon creates 0660 root:packetwyrm in production --
# so this gateway connects via its packetwyrm group membership, not as root. The
# daemon remains the sole authority for the access secret (forwarded verbatim).
User=packetwyrm
Group=packetwyrm

# Hardening.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectControlGroups=true
ProtectKernelTunables=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
# No writable paths needed (self-signed cert is generated in memory).

[Install]
WantedBy=multi-user.target
