# packaging/packetwyrmd.service
[Unit]
Description=PacketWyrm tester daemon
Documentation=https://github.com/kazubu/packetwyrm
After=network.target
Requires=sys-subsystem-net-devices-tun.device

[Service]
Type=simple
# -p 9100 binds the (unauthenticated) Prometheus /metrics exporter to
# 127.0.0.1 by default. To scrape from another host, change it to
# -p 0.0.0.0:9100 (deliberate opt-in) and firewall it.
ExecStart=/usr/local/sbin/packetwyrmd -c /etc/packetwyrm/packetwyrm.yaml \
                                     -s 5000 -p 9100
Restart=on-failure
RestartSec=2s

# This daemon runs as root and needs the FULL root capability set: it mmaps
# the FPGA PCI BAR via /sys/.../resource0 (CAP_SYS_RAWIO), reads/writes PCI
# config + power sysfs (CAP_SYS_ADMIN), creates TAP devices (CAP_NET_ADMIN),
# and drives VFIO. Do NOT clamp CapabilityBoundingSet to CAP_NET_ADMIN --
# that silently drops the capabilities the hardware backend needs, so
# `systemctl start` fails to mmap the BAR / prep PCI even though a manual
# `sudo packetwyrmd` (full caps) works. Drop to User=packetwyrm + a narrow
# CapabilityBoundingSet only once the udev rule (scripts/99-packetwyrm.rules)
# grants unprivileged BAR access on the operator's setup.
User=root

# systemd-managed directories (created on start, removed on stop).
RuntimeDirectory=packetwyrm
RuntimeDirectoryMode=0755
StateDirectory=packetwyrm
LogsDirectory=packetwyrm
ConfigurationDirectory=packetwyrm

# Hardening. Adjust as needed for your deployment; PrivateNetwork
# would break TAP creation, ProtectKernelTunables blocks sysfs PCI
# scans, so neither is enabled here.
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/run/packetwyrm /var/lib/packetwyrm /var/log/packetwyrm
ProtectHome=true
ProtectControlGroups=true
PrivateTmp=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true

[Install]
WantedBy=multi-user.target
